# capstone_project/daemon/logger.timer
[Unit]
Description=Background Logger Timer

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target

