# capstone_project/daemon/logger.service
[Unit]
Description=Background Logger
After=network.target

[Service]
ExecStart=/usr/local/bin/logger

[Install]
WantedBy=multi-user.target
